# systemd/user/waybar.service.d/aptosid.conf
# Launch Waybar with aptosid's per-session config instead of the shared
# /etc/xdg/waybar defaults. All uwsm sessions (Hyprland, Sway, labwc) start
# waybar via this user unit (WantedBy=graphical-session.target); aptosid-waybar
# picks the themed config from $XDG_CURRENT_DESKTOP. Emptying ExecStart first
# is required before redefining it.
[Service]
ExecStart=
ExecStart=/usr/libexec/aptosid-waybar
