# systemd/user/waybar.service.d/aptosid.conf
# Launch Waybar with aptosid's Hyprland config instead of the shared
# /etc/xdg/waybar defaults. Hyprland starts waybar via this user unit
# (WantedBy=graphical-session.target under uwsm); sway/labwc launch it
# themselves. Emptying ExecStart first is required before redefining it.
[Service]
ExecStart=
ExecStart=/usr/libexec/aptosid-waybar hyprland
